# scripts/services/sync.service
[Unit]
Description=Syncing the Codium Configuration Service

[Service]
ExecStart=
Restart=always

[Install]
WantedBy=multi-user.target